# cronjob-hourly.service
[Unit]
Description=Hourly cron jobs

[Service]
ExecStart=/bin/run-parts /etc/cron.hourly

# vi: ft=systemd
