# cronjob-daily.service
[Unit]
Description=Daily cron jobs

[Service]
ExecStart=/bin/run-parts /etc/cron.daily
