# cronjob-monthly.service
[Unit]
Description=Monthly cron jobs

[Service]
ExecStart=/bin/run-parts /etc/cron.monthly

# vi: ft=systemd
